# deploy/systemd/speedtest-next.service
[Unit]
Description=Speedtest Next web service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
DynamicUser=yes
WorkingDirectory=/opt/speedtest-next
ExecStart=/opt/speedtest-next/speedtest-next
Restart=always
RestartSec=3

Environment=SPEEDTEST_ADDR=:8080
Environment=SPEEDTEST_LOG_NOISY_API=false

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true

[Install]
WantedBy=multi-user.target
